# deployments/unified-stock-trader/supervisor.conf
[program:unified-stock-trader]
; Meta selector (2026-03-18): expanded cache-ready live universe
command=/nvme0n1-disk/code/stock-prediction/.venv313/bin/python -u /nvme0n1-disk/code/stock-prediction/unified_hourly_experiment/trade_unified_hourly_meta.py --strategy wd04=/nvme0n1-disk/code/stock-prediction/unified_hourly_experiment/checkpoints/wd_0.04:9 --strategy wd06=/nvme0n1-disk/code/stock-prediction/unified_hourly_experiment/checkpoints/wd_0.06_s42:8 --strategy wd05=/nvme0n1-disk/code/stock-prediction/unified_hourly_experiment/checkpoints/wd_0.05_s42:19 --strategy wd08=/nvme0n1-disk/code/stock-prediction/unified_hourly_experiment/checkpoints/wd_0.08_s42:10 --strategy wd03=/nvme0n1-disk/code/stock-prediction/unified_hourly_experiment/checkpoints/wd_0.03_s42:20 --strategy robb=/nvme0n1-disk/code/stock-prediction/unified_hourly_experiment/checkpoints/stock_sortino_lag_robust_20260219b_fast_rw012_sm003_lagr01 --strategy robc=/nvme0n1-disk/code/stock-prediction/unified_hourly_experiment/checkpoints/stock_sortino_lag_robust_20260219c_l0123_rw010_sm006_seq48 --stock-symbols NVDA,PLTR,GOOG,DBX,TRIP,MTCH,NYT,AAPL,MSFT,META,TSLA,NET,BKNG,EBAY,EXPE --min-edge 0.001 --fee-rate 0.001 --max-positions 5 --max-hold-hours 5 --trade-amount-scale 100.0 --min-buy-amount 2.0 --entry-intensity-power 1.0 --entry-min-intensity-fraction 0.0 --long-intensity-multiplier 1.0 --short-intensity-multiplier 1.5 --meta-metric p10 --meta-lookback-days 14 --meta-selection-mode sticky --meta-switch-margin 0.005 --meta-min-score-gap 0.0 --meta-recency-halflife-days 0.0 --meta-history-days 120 --sit-out-if-negative --sit-out-threshold -0.001 --market-order-entry --bar-margin 0.0005 --entry-order-ttl-hours 6 --margin-rate 0.0625 --live --loop
directory=/nvme0n1-disk/code/stock-prediction
user=administrator
autostart=true
autorestart=true
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/unified-stock-trader.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/unified-stock-trader-error.log
stderr_logfile_maxbytes=10MB
environment=PYTHONPATH="/nvme0n1-disk/code/stock-prediction",PYTHONUNBUFFERED="1",CHRONOS2_FREQUENCY="hourly"

[program:stock-cache-refresh]
command=/nvme0n1-disk/code/stock-prediction/.venv313/bin/python -u /nvme0n1-disk/code/stock-prediction/scripts/refresh_stock_caches.py
directory=/nvme0n1-disk/code/stock-prediction
user=administrator
autostart=true
autorestart=true
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/stock-cache-refresh.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/stock-cache-refresh-error.log
stderr_logfile_maxbytes=5MB
environment=PYTHONPATH="/nvme0n1-disk/code/stock-prediction",PYTHONUNBUFFERED="1",CHRONOS2_FREQUENCY="hourly"
