# deployments/binance-hourly-xgb-margin-pack/supervisor.conf
[program:binance-hourly-xgb-margin-pack]
; Single live Binance writer for the hourly XGB cross-margin portfolio pack.
; It refreshes hourly data, runs one gated execute cycle, then sleeps until the
; next hourly cycle. Live order placement is still gated in Python by
; ALLOW_BINANCE_XGB_LIVE_TRADING, writer isolation, fresh data, and exit coverage.
command=/home/lee/code/stock/deployments/binance-hourly-xgb-margin-pack/launch.sh
environment=PYTHONUNBUFFERED="1",HOME="/home/lee"
directory=/home/lee/code/stock
user=lee
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/binance-hourly-xgb-margin-pack-error.log
stdout_logfile=/var/log/supervisor/binance-hourly-xgb-margin-pack.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
