# systemd/alpaca-hourly-trader.service
[Unit]
Description=Alpaca hourly trader (live)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/nvme0n1-disk/code/stock-prediction
Environment=PYTHONPATH=/nvme0n1-disk/code/stock-prediction
Environment=PYTHONUNBUFFERED=1
Environment=PAPER=0
Environment=TRADE_STATE_SUFFIX=hourly
ExecStart=/nvme0n1-disk/code/stock-prediction/.venv313-3/bin/python -m newnanoalpacahourlyexp.trade_alpaca_hourly \
  --symbols SOLUSD,LINKUSD,UNIUSD \
  --checkpoints "SOLUSD=/nvme0n1-disk/code/stock-prediction/binanceneural/checkpoints/alpaca_solusd_20260204_013044/epoch_006.pt,LINKUSD=/nvme0n1-disk/code/stock-prediction/binanceneural/checkpoints/alpaca_linkusd_20260204_111049/epoch_006.pt,UNIUSD=/nvme0n1-disk/code/stock-prediction/binanceneural/checkpoints/alpaca_uniusd_20260204_231926/epoch_004.pt" \
  --default-checkpoint /nvme0n1-disk/code/stock-prediction/binanceneural/checkpoints/alpaca_solusd_20260204_013044/epoch_006.pt \
  --sequence-length 96 \
  --horizon 1 \
  --forecast-horizons 1,24 \
  --forecast-horizons-map "UNIUSD=1" \
  --forecast-cache-root /nvme0n1-disk/code/stock-prediction/binanceneural/forecast_cache \
  --crypto-data-root /nvme0n1-disk/code/stock-prediction/trainingdatahourly/crypto \
  --stock-data-root /nvme0n1-disk/code/stock-prediction/trainingdatahourly/stocks \
  --allocation-pct 0.05 \
  --intensity-scale 1.0 \
  --min-gap-pct 0.001 \
  --buffer-seconds 30
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
