# systemd/alpaca-hourly-trader.service
[Unit]
Description=Alpaca hourly crypto trader (live) - ETHUSD + BTCUSD best checkpoints
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/nvme0n1-disk/code/stock-prediction
Environment=PYTHONPATH=/nvme0n1-disk/code/stock-prediction
Environment=PYTHONUNBUFFERED=1
Environment=PAPER=0
Environment=TRADE_STATE_SUFFIX=hourly
ExecStart=/nvme0n1-disk/code/stock-prediction/.venv313/bin/python -m newnanoalpacahourlyexp.trade_alpaca_hourly \
  --symbols ETHUSD,BTCUSD \
  --checkpoints ETHUSD=/nvme0n1-disk/code/stock-prediction/binanceneural/checkpoints/seed42_ethusd_ft30_20260209_014309/epoch_028.pt,BTCUSD=/nvme0n1-disk/code/stock-prediction/binanceneural/checkpoints/btcusd_h1only_ft30_20260208/epoch_029.pt \
  --sequence-length 72 \
  --horizon 1 \
  --forecast-horizons 1,24 \
  --forecast-cache-root /nvme0n1-disk/code/stock-prediction/binanceneural/forecast_cache \
  --crypto-data-root /nvme0n1-disk/code/stock-prediction/trainingdatahourly/crypto \
  --stock-data-root /nvme0n1-disk/code/stock-prediction/trainingdatahourly/stocks \
  --allocation-pct 0.5 \
  --allocation-mode portfolio \
  --intensity-scale 1.0 \
  --min-gap-pct 0.001 \
  --moving-average-windows 24,72 \
  --ema-windows 24,72 \
  --atr-windows 24,72 \
  --trend-windows 72 \
  --drawdown-windows 72 \
  --volume-z-window 72 \
  --volume-shock-window 24 \
  --vol-regime-short 24 \
  --vol-regime-long 72 \
  --min-history-hours 50 \
  --buffer-seconds 30
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
