# systemd/daily-rl-trader.service
[Unit]
Description=Daily Stock RL Trading Bot
After=network.target

[Service]
Type=simple
User=administrator
WorkingDirectory=/nvme0n1-disk/code/stock-prediction
ExecStart=/bin/bash -lc 'export HOME=/home/administrator && export USER=administrator && cd /nvme0n1-disk/code/stock-prediction && exec deployments/daily-rl-trader/launch.sh'
Restart=on-failure
RestartSec=60
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
