# systemd/daily-rl-trader.service
[Unit]
Description=Daily Mixed23 RL+Gemini Trading Bot
After=network.target

[Service]
Type=simple
User=administrator
WorkingDirectory=/nvme0n1-disk/code/stock-prediction
ExecStart=/nvme0n1-disk/code/stock-prediction/.venv313/bin/python -u trade_mixed_daily.py --daemon --mode hybrid --paper --top-k 2 --max-positions 2 --max-position-fraction 0.25 --max-hold-days 5 --llm-model gemini-3.1-flash-lite-preview --thinking-level=
Restart=on-failure
RestartSec=60
Environment=PYTHONPATH=/nvme0n1-disk/code/stock-prediction

[Install]
WantedBy=multi-user.target
