# supervisor/autoresearch-stock-agent.conf
[program:autoresearch-stock-agent]
command=/nvme0n1-disk/code/stock-prediction/.venv312/bin/python -u -m autoresearch_stock.agent_scheduler --analysis-dir /nvme0n1-disk/code/stock-prediction/analysis/autoresearch_stock_agent --experiment-bundle-root /nvme0n1-disk/code/stock-prediction/experiments/autoresearch_stock_agent --repo-root /nvme0n1-disk/code/stock-prediction --python /nvme0n1-disk/code/stock-prediction/.venv312/bin/python --backends codex --frequencies hourly,daily --loop --max-turns 1000000 --interval-seconds 300 --codex-reasoning-effort xhigh
directory=/nvme0n1-disk/code/stock-prediction
autostart=false
autorestart=false
stopwaitsecs=120
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/autoresearch-stock-agent.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/autoresearch-stock-agent-error.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
environment=PYTHONPATH="/nvme0n1-disk/code/stock-prediction",PYTHONUNBUFFERED="1"
