# supervisor/binance-5min-collector.conf
[program:binance-5min-collector]
command=/home/lee/code/stock/.venv313/bin/python -u /home/lee/code/stock/scripts/collect_5min_bars.py --symbols DOGEUSDT
environment=PYTHONUNBUFFERED="1",HOME="/home/lee",PYTHONPATH="/home/lee/code/stock"
directory=/home/lee/code/stock
user=lee
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/binance-5min-collector-error.log
stdout_logfile=/var/log/supervisor/binance-5min-collector.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
