# supervisor/binance-doge-margin.conf
[program:binance-doge-margin]
; DOGE 1x margin trading bot - R4_h384_cosine ep1 (safety=8.59, sort=39.9, ret=+94%, dd=-3.1%)
command=/home/lee/code/stock/.venv313/bin/python -u -m binanceleveragesui.trade_margin_sui --symbol DOGEUSDT --data-symbol DOGEUSD --checkpoint /home/lee/code/stock/binanceleveragesui/checkpoints/DOGEUSD_r4_R4_h384_cosine/binanceneural_20260301_065549/epoch_001.pt --max-leverage 1.0 --sequence-length 72 --horizon 1 --intensity-scale 5.0 --min-gap-pct 0.0003 --max-hold-hours 6 --poll-seconds 30 --expiry-minutes 90 --cycle-minutes 5 --data-root trainingdatahourlybinance --forecast-cache binanceneural/forecast_cache --forecast-horizons 1
environment=HF_HOME="/vfast/data/stock/hf_cache",HUGGINGFACE_HUB_CACHE="/vfast/data/stock/hf_cache",TRANSFORMERS_CACHE="/vfast/data/stock/hf_cache",CHRONOS2_MODEL_ID_OVERRIDE="/vfast/data/stock/hf_cache/models--amazon--chronos-2/snapshots/0f8a440441931157957e2be1a9bce66627d99c76",PAPER="0",PYTHONUNBUFFERED="1",HOME="/home/lee",PYTHONPATH="/home/lee/code/stock"
directory=/home/lee/code/stock
user=lee
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/binance-doge-margin-error.log
stdout_logfile=/var/log/supervisor/binance-doge-margin.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
