# supervisor/binanceexp1-cache-refresh.conf
[program:binanceexp1-cache-refresh]
command=/home/lee/code/stock/.venv313/bin/python -u /home/lee/code/stock/scripts/refresh_binanceexp1_caches.py
environment=HF_HOME="/vfast/data/stock/hf_cache",HUGGINGFACE_HUB_CACHE="/vfast/data/stock/hf_cache",TRANSFORMERS_CACHE="/vfast/data/stock/hf_cache",CHRONOS2_MODEL_ID_OVERRIDE="/vfast/data/stock/hf_cache/models--amazon--chronos-2/snapshots/0f8a440441931157957e2be1a9bce66627d99c76",PYTHONUNBUFFERED="1",HOME="/home/lee",PYTHONPATH="/home/lee/code/stock"
directory=/home/lee/code/stock
user=lee
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/binanceexp1-cache-refresh-error.log
stdout_logfile=/var/log/supervisor/binanceexp1-cache-refresh.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
