# supervisor/binance-sui.conf
[program:binance-sui]
; Updated 2026-02-14: Best (sui_sortino_rw0012_lr1e4) - 318.7% return, Sortino 786
command=/home/lee/code/stock/.venv313/bin/python -u -m bitbankstylelongsuitrain.trade_binance_sui --checkpoint /home/lee/code/stock/binancechronossolexperiment/checkpoints/sui_sortino_rw0012_lr1e4/policy_checkpoint.pt --sequence-length 72 --horizon 1 --intensity-scale 6.0 --price-offset-pct 0.0005 --min-gap-pct 0.0003 --poll-seconds 30 --expiry-minutes 90 --cycle-minutes 5 --refresh-data
environment=HF_HOME="/vfast/data/stock/hf_cache",HUGGINGFACE_HUB_CACHE="/vfast/data/stock/hf_cache",TRANSFORMERS_CACHE="/vfast/data/stock/hf_cache",CHRONOS2_MODEL_ID_OVERRIDE="/vfast/data/stock/hf_cache/models--amazon--chronos-2/snapshots/0f8a440441931157957e2be1a9bce66627d99c76",PAPER="0",PYTHONUNBUFFERED="1",HOME="/home/lee",PYTHONPATH="/home/lee/code/stock",BINANCE_DEFAULT_QUOTE="FDUSD"
directory=/home/lee/code/stock
user=lee
autostart=false
autorestart=true
stderr_logfile=/var/log/supervisor/binance-sui-error.log
stdout_logfile=/var/log/supervisor/binance-sui.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
