# supervisor/binance-meta-margin.conf
[program:binance-meta-margin]
; Frozen deploy promoted on 2026-03-08 after live frontier validation.
; Strategy unchanged: DOGE/AAVE epoch_003 pair with winner_cash + omega, lookback=1,
; allow_short, max_short_leverage=0.16, profit gate 24h > 0, no probe cap.
; v2 rebundles the live app copy with dust-sizing and replay-parity fixes.
command=/home/lee/code/stock/deployments/binance-meta-margin/20260308_h1h6_e3_omega_s016_gate24_full_v2/launch_frozen.sh
environment=HF_HOME="/vfast/data/stock/hf_cache",HUGGINGFACE_HUB_CACHE="/vfast/data/stock/hf_cache",TRANSFORMERS_CACHE="/vfast/data/stock/hf_cache",CHRONOS2_MODEL_ID_OVERRIDE="/vfast/data/stock/hf_cache/models--amazon--chronos-2/snapshots/0f8a440441931157957e2be1a9bce66627d99c76",PAPER="0",PYTHONUNBUFFERED="1",HOME="/home/lee"
directory=/home/lee/code/stock/deployments/binance-meta-margin/20260308_h1h6_e3_omega_s016_gate24_full_v2
user=lee
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/binance-meta-margin-error.log
stdout_logfile=/var/log/supervisor/binance-meta-margin.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
