# supervisor/binance-5min-collector.conf
[program:binance-5min-collector]
command=/nvme0n1-disk/code/stock-prediction/.venv/bin/python -u /nvme0n1-disk/code/stock-prediction/scripts/collect_5min_bars.py --symbols DOGEUSDT BTCUSDT ETHUSDT SOLUSDT BNBUSDT XRPUSDT LINKUSDT AVAXUSDT DOTUSDT AAVEUSDT UNIUSDT LTCUSDT ALGOUSDT OPUSDT ARBUSDT MATICUSDT SHIBUSDT
environment=PYTHONUNBUFFERED="1",HOME="/home/administrator",PYTHONPATH="/nvme0n1-disk/code/stock-prediction"
directory=/nvme0n1-disk/code/stock-prediction
user=administrator
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/binance-5min-collector-error.log
stdout_logfile=/var/log/supervisor/binance-5min-collector.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
