# supervisor/binance-doge-margin.conf
[program:binance-doge-margin]
; Meta margin bot (DOGE + AAVE) selected from epoch pair sweep.
; Probe deployment only until replayability is trustworthy again.
; Pair: doge_drop15_ep008 + aave_strides_ep002
; max-position-notional=5.10 keeps live trades near 5 USDT while still clearing exchange min-notional/lot rounding.
command=/home/lee/code/stock/.venv313/bin/python -u -m binanceleveragesui.trade_margin_meta --doge-checkpoint /home/lee/code/stock/binanceleveragesui/checkpoints/r5_DOGE_rw05_drop15/binanceneural_20260303_001154/epoch_008.pt --aave-checkpoint /home/lee/code/stock/binanceleveragesui/checkpoints/r5_AAVE_rw05_strides_long/binanceneural_20260303_023521/epoch_002.pt --max-leverage 2.30 --max-position-notional 5.10 --lookback 24 --selection-mode winner_cash --selection-metric calmar --cash-threshold 0.0 --switch-margin 0.005 --sequence-length 72 --horizon 1 --intensity-scale 5.0 --min-gap-pct 0.0003 --max-hold-hours 6 --poll-seconds 30 --expiry-minutes 90 --cycle-minutes 60 --data-root trainingdatahourlybinance --forecast-cache binanceneural/forecast_cache --forecast-horizons 1
environment=HF_HOME="/vfast/data/stock/hf_cache",HUGGINGFACE_HUB_CACHE="/vfast/data/stock/hf_cache",TRANSFORMERS_CACHE="/vfast/data/stock/hf_cache",CHRONOS2_MODEL_ID_OVERRIDE="/vfast/data/stock/hf_cache/models--amazon--chronos-2/snapshots/0f8a440441931157957e2be1a9bce66627d99c76",PAPER="0",PYTHONUNBUFFERED="1",HOME="/home/lee",PYTHONPATH="/home/lee/code/stock"
directory=/home/lee/code/stock
user=lee
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/binance-doge-margin-error.log
stdout_logfile=/var/log/supervisor/binance-doge-margin.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
