# supervisor/pufferlib-stocks-live.conf
[program:pufferlib-stocks-live]
command=/home/lee/code/stock/.venv/bin/python -m pufferlib_market.trade_ppo_stocks --checkpoint experiments/pufferlib_stocks7_50M/best.pt --allocation-usd 1000
directory=/home/lee/code/stock
user=lee
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/pufferlib-stocks-live.log
stderr_logfile=/var/log/supervisor/pufferlib-stocks-live-error.log
environment=PAPER="0",PYTHONPATH="/home/lee/code/stock"
