# supervisor/pufferlib-stocks-paper.conf
[program:pufferlib-stocks-paper]
command=/home/lee/code/stock/.venv/bin/python -m pufferlib_market.trade_ppo_stocks --checkpoint experiments/pufferlib_stocks7_50M/best.pt --paper --allocation-usd 1000
directory=/home/lee/code/stock
user=lee
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/pufferlib-stocks-paper.log
stderr_logfile=/var/log/supervisor/pufferlib-stocks-paper-error.log
environment=PAPER="1",PYTHONPATH="/home/lee/code/stock"
