# supervisor/stock-5min-collector.local-supervisord.conf
[unix_http_server]
file=/nvme0n1-disk/code/stock-prediction/strategy_state/stock-5min-supervisord.sock
chmod=0700

[supervisord]
logfile=/nvme0n1-disk/code/stock-prediction/strategy_state/stock-5min-supervisord.log
pidfile=/nvme0n1-disk/code/stock-prediction/strategy_state/stock-5min-supervisord.pid
childlogdir=/nvme0n1-disk/code/stock-prediction/strategy_state
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///nvme0n1-disk/code/stock-prediction/strategy_state/stock-5min-supervisord.sock

[program:stock-5min-collector]
command=/nvme0n1-disk/code/stock-prediction/.venv313/bin/python -u /nvme0n1-disk/code/stock-prediction/scripts/collect_stock_5min_bars.py --supervisor-config /nvme0n1-disk/code/stock-prediction/supervisor/unified-stock-trader.conf --extra-symbols TSLA --out-root /nvme0n1-disk/code/stock-prediction/trainingdata5min --interval-seconds 300 --lookback-minutes 15 --overlap-minutes 10 --bootstrap-days 30 --repair-days 7 --repair-interval-hours 24
environment=PYTHONUNBUFFERED="1",HOME="/home/administrator",PYTHONPATH="/nvme0n1-disk/code/stock-prediction"
directory=/nvme0n1-disk/code/stock-prediction
autostart=true
autorestart=true
stdout_logfile=/nvme0n1-disk/code/stock-prediction/strategy_state/stock-5min-collector.log
stderr_logfile=/nvme0n1-disk/code/stock-prediction/strategy_state/stock-5min-collector-error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
stopasgroup=true
killasgroup=true
