# supervisor/stock-5min-collector.conf
[program:stock-5min-collector]
command=/nvme0n1-disk/code/stock-prediction/.venv313/bin/python -u /nvme0n1-disk/code/stock-prediction/scripts/collect_stock_5min_bars.py --supervisor-config /nvme0n1-disk/code/stock-prediction/supervisor/unified-stock-trader.conf --extra-symbols TSLA --out-root /nvme0n1-disk/code/stock-prediction/trainingdata5min --interval-seconds 300 --lookback-minutes 15 --overlap-minutes 10 --bootstrap-days 30 --repair-days 7 --repair-interval-hours 24
environment=PYTHONUNBUFFERED="1",HOME="/home/administrator",PYTHONPATH="/nvme0n1-disk/code/stock-prediction"
directory=/nvme0n1-disk/code/stock-prediction
user=administrator
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/stock-5min-collector-error.log
stdout_logfile=/var/log/supervisor/stock-5min-collector.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
