# supervisor/unified-stock-trader-meta.conf
[program:unified-stock-trader-meta]
; Meta-selector stock trader (5-strategy ensemble, sit-out aware)
; Best 2026-03-04 robust sweep:
; metric=sharpe lookback=16 mode=winner sit_out_threshold=0.3 min_edge=0.0065 hold=5
command=/home/lee/code/stock/.venv313/bin/python -u /home/lee/code/stock/unified_hourly_experiment/trade_unified_hourly_meta.py --strategy wd04=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.04:9 --strategy wd06=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.06_s42:8 --strategy wd05=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.05_s42:19 --strategy wd08=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.08_s42:10 --strategy wd03=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.03_s42:20 --stock-symbols NVDA,PLTR,GOOG,DBX,TRIP,MTCH --min-edge 0.0065 --fee-rate 0.001 --max-positions 5 --max-hold-hours 5 --meta-metric sharpe --meta-lookback-days 16 --meta-selection-mode winner --meta-switch-margin 0.0 --meta-min-score-gap 0.0 --meta-history-days 120 --sit-out-if-negative --sit-out-threshold 0.3 --bar-margin 0.0013 --live --loop
directory=/home/lee/code/stock
user=lee
autostart=false
autorestart=true
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/unified-stock-trader-meta.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/unified-stock-trader-meta-error.log
stderr_logfile_maxbytes=10MB
environment=PYTHONPATH="/home/lee/code/stock",PYTHONUNBUFFERED="1"
