# supervisor/unified-stock-trader.conf
[program:unified-stock-trader]
command=/home/lee/code/stock/.venv/bin/python /home/lee/code/stock/unified_hourly_experiment/trade_unified_hourly.py --checkpoint-dir /home/lee/code/stock/unified_hourly_experiment/checkpoints/binanceneural_20260213_181304 --stock-symbols NVDA,MSFT,META,GOOG --dry-run
directory=/home/lee/code/stock
user=lee
autostart=false
autorestart=false
stopwaitsecs=30
stdout_logfile=/home/lee/code/stock/logs/unified-stock-trader.log
stderr_logfile=/home/lee/code/stock/logs/unified-stock-trader.err
environment=PYTHONPATH="/home/lee/code/stock"
