# supervisor/unified-stock-trader.conf
[program:unified-stock-trader]
; fnoise02 ep2 5-stock (NVDA,GOOG,DBX long + TRIP,MTCH short): 7d +0.03%, 14d +0.92%, 30d +1.34%, 60d +1.49%, 90d +1.55%
command=/home/lee/code/stock/.venv313/bin/python -u /home/lee/code/stock/unified_hourly_experiment/trade_unified_hourly.py --checkpoint-dir /home/lee/code/stock/unified_hourly_experiment/checkpoints/top9_lag1_fnoise02 --epoch 2 --stock-symbols NVDA,GOOG,DBX,TRIP,MTCH --min-edge 0.0 --fee-rate 0.001 --max-positions 5 --max-hold-hours 6 --live --loop
directory=/home/lee/code/stock
user=lee
autostart=true
autorestart=true
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/unified-stock-trader.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/unified-stock-trader-error.log
stderr_logfile_maxbytes=10MB
environment=PYTHONPATH="/home/lee/code/stock",PYTHONUNBUFFERED="1"

[program:stock-cache-refresh]
command=/home/lee/code/stock/.venv313/bin/python -u /home/lee/code/stock/scripts/refresh_stock_caches.py
directory=/home/lee/code/stock
user=lee
autostart=true
autorestart=true
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/stock-cache-refresh.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/stock-cache-refresh-error.log
stderr_logfile_maxbytes=5MB
environment=PYTHONPATH="/home/lee/code/stock",PYTHONUNBUFFERED="1"
