# supervisor/unified-stock-trader.conf
[program:unified-stock-trader]
; Meta selector (2026-03-08): limit-entry first_trigger refine; sticky p10, lookback=18, switch=0.005, sit_out=-0.001, edge=0.001, min_buy=2.0
; Auto-updated 2026-03-13: min_goodness=2.9381 mean_goodness=3.3502 min_ret=0.2976%
command=/home/lee/code/stock/.venv313/bin/python -u /home/lee/code/stock/unified_hourly_experiment/trade_unified_hourly_meta.py --strategy wd04=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.04:9 --strategy wd06=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.06_s42:8 --strategy wd05=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.05_s42:19 --strategy wd08=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.08_s42:10 --strategy wd03=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.03_s42:20 --stock-symbols NVDA,PLTR,GOOG,DBX,TRIP,MTCH,AAPL --min-edge 0.001 --fee-rate 0.001 --max-positions 5 --max-hold-hours 5 --trade-amount-scale 100.0 --min-buy-amount 2.0 --entry-intensity-power 1.0 --entry-min-intensity-fraction 0.0 --long-intensity-multiplier 1.0 --short-intensity-multiplier 1.5 --meta-metric p10 --meta-lookback-days 18 --meta-selection-mode sticky --meta-switch-margin 0.005 --meta-min-score-gap 0.0 --meta-recency-halflife-days 0.0 --meta-history-days 120 --sit-out-if-negative --sit-out-threshold -0.001 --entry-selection-mode edge_rank --bar-margin 0.0005 --entry-order-ttl-hours 0 --margin-rate 0.0625 --live --loop --market-order-entry
directory=/home/lee/code/stock
user=lee
autostart=false
autorestart=false
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/unified-stock-trader.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/unified-stock-trader-error.log
stderr_logfile_maxbytes=10MB
environment=PYTHONPATH="/home/lee/code/stock",PYTHONUNBUFFERED="1",CHRONOS2_FREQUENCY="hourly"

[program:stock-deleverage-bot]
; Runs during final 60min before market close, checks every 5min,
; sells positions via limit orders (0.05% below market) to keep leverage <= 2x overnight
command=/home/lee/code/stock/.venv313/bin/python -u -m unified_orchestrator.deleverage_bot --live --deleverage-window 60 --check-interval 300 --target 2.0
directory=/home/lee/code/stock
user=lee
autostart=true
autorestart=true
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/stock-deleverage-bot.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/var/log/supervisor/stock-deleverage-bot-error.log
stderr_logfile_maxbytes=5MB
environment=PYTHONPATH="/home/lee/code/stock",PYTHONUNBUFFERED="1"

[program:stock-cache-refresh]
; Auto-updated 2026-03-13: min_goodness=2.9381 mean_goodness=3.3502 min_ret=0.2976%
command=/home/lee/code/stock/.venv313/bin/python -u /home/lee/code/stock/unified_hourly_experiment/trade_unified_hourly_meta.py --strategy wd04=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.04:9 --strategy wd06=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.06_s42:8 --strategy wd05=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.05_s42:19 --strategy wd08=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.08_s42:10 --strategy wd03=/home/lee/code/stock/unified_hourly_experiment/checkpoints/wd_0.03_s42:20 --stock-symbols NVDA,PLTR,GOOG,DBX,TRIP,MTCH,AAPL --min-edge 0.001 --fee-rate 0.001 --max-positions 5 --max-hold-hours 5 --trade-amount-scale 100.0 --min-buy-amount 2.0 --entry-intensity-power 1.0 --entry-min-intensity-fraction 0.0 --long-intensity-multiplier 1.0 --short-intensity-multiplier 1.5 --meta-metric p10 --meta-lookback-days 18 --meta-selection-mode sticky --meta-switch-margin 0.005 --meta-min-score-gap 0.0 --meta-recency-halflife-days 0.0 --meta-history-days 120 --sit-out-if-negative --sit-out-threshold -0.001 --entry-selection-mode edge_rank --bar-margin 0.0005 --entry-order-ttl-hours 0 --margin-rate 0.0625 --live --loop --market-order-entry
directory=/home/lee/code/stock
user=lee
autostart=false
autorestart=false
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/stock-cache-refresh.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/stock-cache-refresh-error.log
stderr_logfile_maxbytes=5MB
environment=PYTHONPATH="/home/lee/code/stock",PYTHONUNBUFFERED="1",CHRONOS2_FREQUENCY="hourly"
