# supervisor/unified-hourly-paper.conf
[program:unified-hourly-paper]
command=/home/lee/code/stock/.venv/bin/python unified_hourly_experiment/trade_unified_hourly.py --checkpoint-dir unified_hourly_experiment/checkpoints/unified_v3_moredata --stock-symbols "" --crypto-symbols ETHUSD,SOLUSD,UNIUSD --crypto-cache-root binanceneural/forecast_cache --paper --loop
directory=/home/lee/code/stock
user=lee
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/unified-hourly-paper.log
stderr_logfile=/var/log/supervisor/unified-hourly-paper-error.log
environment=PAPER="1"
