# supervisor/unified-stock-trader.conf
[program:unified-stock-trader]
; lag=1 lr=3e-5 ep3 (sortino 2.15, ret +1.94% on 30d holdout)
command=/home/lee/code/stock/.venv313/bin/python -u /home/lee/code/stock/unified_hourly_experiment/trade_unified_hourly.py --checkpoint-dir /home/lee/code/stock/unified_hourly_experiment/checkpoints/top10_h1only_6L_lev2x_lag1_lr3e5 --epoch 3 --stock-symbols NVDA,PLTR,TRIP,EBAY,GOOG,NET,KIND,MTCH,NYT,DBX --min-edge 0.0 --fee-rate 0.001 --max-positions 10 --max-hold-hours 6 --live --loop
directory=/home/lee/code/stock
user=lee
autostart=true
autorestart=true
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/unified-stock-trader.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/unified-stock-trader-error.log
stderr_logfile_maxbytes=10MB
environment=PYTHONPATH="/home/lee/code/stock",PYTHONUNBUFFERED="1"

[program:stock-cache-refresh]
command=/home/lee/code/stock/.venv313/bin/python -u /home/lee/code/stock/scripts/refresh_stock_caches.py
directory=/home/lee/code/stock
user=lee
autostart=true
autorestart=true
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/stock-cache-refresh.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/stock-cache-refresh-error.log
stderr_logfile_maxbytes=5MB
environment=PYTHONPATH="/home/lee/code/stock",PYTHONUNBUFFERED="1"
