# deployments/daily-stock-ppo/supervisor.conf
[program:daily-stock-ppo]
; Ensemble PPO daily stock trader: tp05_s123 + tp05_s15 (h=1024, trade_penalty=0.05, softmax_avg)
; Val perf: med=+28.76%/90d, p10=+16.37%, worst=+10.17%, 0/50 neg @ 5bps (up from med=16.52% standalone)
; Trades 12 symbols: AAPL MSFT NVDA GOOG META TSLA SPY QQQ PLTR JPM V AMZN
; NOTE: unified-stock-trader also trades AAPL/MSFT/NVDA/GOOG/META/TSLA/PLTR
; WARNING: symbol conflicts with unified-stock-trader — coordinate before enabling both
command=/nvme0n1-disk/code/stock-prediction/.venv313/bin/python -u \
    /nvme0n1-disk/code/stock-prediction/trade_daily_stock_prod.py \
    --daemon \
    --live \
    --data-source alpaca \
    --allocation-pct 25.0
directory=/nvme0n1-disk/code/stock-prediction
user=administrator
autostart=false
autorestart=true
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/daily-stock-ppo.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/daily-stock-ppo-error.log
stderr_logfile_maxbytes=10MB
environment=PYTHONPATH="/nvme0n1-disk/code/stock-prediction",PYTHONUNBUFFERED="1"
