# deployments/binance-worksteal-daily/supervisor.conf
[program:binance-worksteal-daily]
; Work-stealing daily dip-buying strategy across 30 crypto symbols
; Runs at UTC midnight, places limit buy orders 20% below recent highs
; Best 30d: +9.13%, Sort=23.03, -1.5% MaxDD, 73% WR
command=/home/lee/code/stock/deployments/binance-worksteal-daily/launch.sh
environment=PYTHONUNBUFFERED="1",HOME="/home/lee"
directory=/home/lee/code/stock
user=lee
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/binance-worksteal-daily-error.log
stdout_logfile=/var/log/supervisor/binance-worksteal-daily.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
