# old_prod/2026-04-22-crypto-weekend-live-standalone-superseded/supervisor.conf
[program:crypto-weekend-live]
directory=/nvme0n1-disk/code/stock-prediction
command=/bin/bash -lc 'cd /nvme0n1-disk/code/stock-prediction && exec deployments/crypto-weekend-live/launch.sh'
user=administrator
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/crypto-weekend-live.log
stderr_logfile=/var/log/supervisor/crypto-weekend-live-error.log
environment=PYTHONUNBUFFERED="1"
